# contrib/debian/DEBIAN/prerm
#!/bin/sh

set -e

case "$1" in
    remove)
        invoke-rc.d jerrybuild stop
        ;;
esac

exit 0
